# firefly-iii/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Firefly III
# Runs the PHP-FPM daemon
# ==============================================================================

declare key

key=$(cat /data/firefly/appkey.txt)
export APP_KEY=${key}
export LOG_CHANNEL=stdout
export APP_ENV=local
export DB_CONNECTION=mysql
export DB_DATABASE
export DB_HOST
export DB_PASSWORD
export DB_PORT
export DB_USERNAME
export MAPBOX_API_KEY

if bashio::config.has_value 'mapbox_api_key';then
    MAPBOX_API_KEY=$(bashio::config "mapbox_api_key")
fi

DB_HOST=$(bashio::services "mysql" "host")
DB_DATABASE=firefly
DB_USERNAME=$(bashio::services "mysql" "username")
DB_PASSWORD=$(bashio::services "mysql" "password")
DB_PORT=$(bashio::services "mysql" "port")

if bashio::config.equals 'log_level' 'debug' \
    || bashio::config.equals 'log_level' 'trace' ;then
    export APP_DEBUG=true
    export APP_LOG_LEVEL=debug
fi
bashio::log.info "Installing/updating Database"

bashio::log.info "Step 1: migrate:refresh"
php /var/www/firefly/artisan migrate:refresh --seed

bashio::log.info "Step 2: firefly-iii:upgrade-database"
php /var/www/firefly/artisan firefly-iii:upgrade-database

bashio::log.info "Step 3: passport:install"
php /var/www/firefly/artisan passport:install

bashio::log.info "Step 4: key:generate"
php /var/www/firefly/artisan key:generate


#sphp /var/www/firefly/artisan firefly-iii:create-database > /dev/null 2>&1

bashio::log.info "Starting PHP-FPM..."

exec php-fpm8 --nodaemonize
